# egardia_alarm/rootfs/etc/services.d/egardia/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Egardia Alarm Addon..."
exec python3 /usr/bin/egardia_addon.py
